# system/sync-samples/sync-samples.service
[Unit]
Description=Sync samples from external drive to NAS
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sync-samples.sh
Nice=10
IOSchedulingClass=idle
